# systemd/ass.service
[Unit]
Description=Local Assistant CLI
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/USER/ass
ExecStart=/home/USER/ass/.venv/bin/python -m app.chat
# Load secrets from a root-owned env file (not committed to git)
EnvironmentFile=/etc/ass.env
# Optional: SMTP/IMAP for email (set in /etc/ass.env)
# SMTP_HOST=...
# SMTP_USER=...
# SMTP_PASS=...
# SMTP_FROM=...
# IMAP_HOST=...
# IMAP_USER=...
# IMAP_PASS=...
Restart=on-failure
User=USER
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
